# skelly_ui/etc/services.d/skelly/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# hand off to your start script (keeps PID/signals correct under s6)
exec /run.sh
